# laptop/systemd/openwebui/container-ollama.service
# container-ollama.service
# autogenerated by Podman 4.3.1
# Thu Jan  1 23:16:46 PST 2026

[Unit]
Description=Podman container-ollama.service
Documentation=man:podman-generate-systemd(1)
Wants=network-online.target
After=network-online.target
RequiresMountsFor=%t/containers

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
TimeoutStopSec=70
ExecStartPre=/bin/rm \
	-f %t/%n.ctr-id
ExecStart=/usr/bin/podman run \
	--cidfile=%t/%n.ctr-id \
	--cgroups=no-conmon \
	--rm \
	--sdnotify=conmon \
	-d \
	--replace \
	--name ollama \
	-e OLLAMA_CONTEXT_LENGTH=130000 \
	-p 11434:11434 \
	--cpus=7 \
	--memory=15g \
	--memory-reservation=10g \
	-v ollama-data:/root/.ollama:rw docker.io/ollama/ollama:0.13.5
ExecStop=/usr/bin/podman stop \
	--ignore -t 10 \
	--cidfile=%t/%n.ctr-id
ExecStopPost=/usr/bin/podman rm \
	-f \
	--ignore -t 10 \
	--cidfile=%t/%n.ctr-id
Type=notify
NotifyAccess=all

[Install]
WantedBy=default.target
